# etc/profile.d/20-local-dns-resolver.sh
/home/werner/software/anti-gfw/dns/dns.bash